# site/www/docker_build/host_controller/etc/supervisor/conf.d/supervisord-redis.conf
[program:redis]
command=/usr/bin/redis-server --port 6379
directory=/srv/home/var/
user=www-data
stopasgroup = true
killasgroup = true
