# site/www/docker_build/host_controller/etc/supervisor/conf.d/supervisord-redis.conf
[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
directory=/srv/home/var/
user=www-data
stopasgroup = true
killasgroup = true
