# docker/nodejs/supervisord-nodejs.conf
[program:webserver]
command=/usr/local/start_app
directory=/srv/home/www/
user=www-data
