# docker/nodejs/supervisord-nodejs.conf
[program:webserver]
command=/usr/local/bin/npm start
directory=/srv/home/www/
user=www-data
