# app_containers/python/etc/supervisor/conf.d/supervisord-python.conf
[program:web]
command=/usr/local/start_app
directory=/srv/home/www/
user=www-data
stopasgroup = true
killasgroup = true
