# site/www/docker_build/host_worker/etc/systemd/system/opendkim.service.d/gen_key.conf
[Service]
ExecStartPre=/usr/local/generate_opendkdim_key
