# docker/base/etc/systemd/system/ssh.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/sshd -D -f /srv/home/etc/ssh/sshd_config

