# site/www/docker_build/host_worker/etc/supervisor/conf.d/supervisord-celery.conf
[program:celery]
command=pyenv /srv/home/www/start_host_worker
directory=/srv/home/var/
user=root
autorestart=false
startsecs=0
