# site/www/docker_build/host_worker/etc/supervisor/conf.d/supervisord-domains.conf
[program:domains]
command=pyenv /srv/home/www/start_domains
directory=/srv/home/var/
autorestart=false
